# contrib/dehydrated@.service.d/reload-nginx.conf
[Service]
ExecStartPost=!-systemctl --no-block reload nginx
